# discodj-debian/discodj-debian-war/src/deb/control/postinst
#!/bin/sh

update-rc.d discodj defaults
service discodj start
