# discodj-debian/discodj-debian-war/src/deb/control/prerm
#!/bin/sh

update-rc.d -f discodj remove
service discodj stop
